# recipes/base-files/files/watchdog
#!/bin/sh /etc/rc.common 
# Copyright (C) 2008 OpenWrt.org

START=97
start() {
	insmod /lib/modules/2.6.35/net/watchdog.ko
}
